# deploy/systemd/ai-diviner.service
[Unit]
Description=AI Diviner Node Service
After=network.target

[Service]
Type=simple
User=__APP_USER__
WorkingDirectory=__APP_DIR__
Environment=NODE_ENV=production
EnvironmentFile=__APP_DIR__/.env
ExecStart=/usr/bin/node server/index.js
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target